# bin/node.conf
[program:zbclient]
user=root
command=/usr/bin/zbclient
autostart = true
autorestart = true

